# deploy/systemd/hhru-daily-backup.service
[Unit]
Description=HH.ru platform daily PostgreSQL backup pipeline
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
OnFailure=hhru-ops-failure-notify@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/hh_collector
Environment=HHRU_BACKUP_DAILY_LOCAL_RETENTION_DAYS=2
EnvironmentFile=-/etc/hhru-platform/backup-daily.env
ExecStart=/opt/hh_collector/scripts/ops/run_daily_backup.sh
TimeoutStartSec=infinity
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
