# deploy/systemd/hhru-weekly-backup-restore-drill.service
[Unit]
Description=HH.ru platform weekly offsite backup restore drill
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
OnFailure=hhru-ops-failure-notify@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/hh_collector
EnvironmentFile=-/etc/hhru-platform/backup-restore-drill.env
ExecStart=/opt/hh_collector/scripts/ops/run_weekly_backup_restore_drill.sh
TimeoutStartSec=infinity
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
