# deploy/systemd/hhru-research-archive.timer
[Unit]
Description=Run HH.ru platform research archive pipeline daily

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=hhru-research-archive.service

[Install]
WantedBy=timers.target
