# deploy/systemd/hhru-research-archive.service
[Unit]
Description=HH.ru platform daily research archive pipeline
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hh_collector
EnvironmentFile=-/etc/hhru-platform/research-archive-daily.env
ExecStart=/opt/hh_collector/scripts/ops/run_daily_research_archive.sh
TimeoutStartSec=infinity
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
